# packages/courDeCassation/LabelServerProd.service
[Unit]
Description=Label server
After=syslog.target

[Service]
User=label
Type=simple
ExecStart=/usr/bin/yarn startProd
WorkingDirectory=/home/label/label
Restart=always
RestartSec=10
KillMode=process

[Install]
WantedBy=multi-user.target
