# packages/courDeCassation/LabelServerPreProd.service
[Unit]
Description=Label server
After=syslog.target

[Service]
Type=simple
ExecStart=/opt/yarn-v1.22.5/bin/yarn startPreProd
WorkingDirectory=/home/label/git/label
Restart=always
RestartSec=10
KillMode=process

[Install]
WantedBy=multi-user.target
